# scripts/coord/com.nobulart.octanex.coord-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nobulart.octanex.coord-watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/craig/octanex-mcp/scripts/coord/watch.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/craig/octanex-mcp</string>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>StandardOutPath</key>
  <string>/Users/craig/hermes-bridge/coord.cron.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/craig/hermes-bridge/coord.cron.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
